# utils/nightly/getdp_linux32_nightly.cron
#
# this will launch a new build every night at 1h31
#
31 1 * * * cd ${HOME}/src/getdp/utils/nightly && /usr/local/bin/ctest -S getdp_linux32_nightly.ctest -VV > ${HOME}/src/getdp/nightly.log 2>&1
