# utils/nightly/mac32_getdp.crontab
#
# this will launch a new build every night at 0h31
#
31 0 * * * ${HOME}/src/getdp/utils/nightly/mac32.sh
